# systemd/user/lib.service
[Unit]
Description=Lib check

[Service]
Type=oneshot
ExecStart=%h/.local/bin/lib -q -s
ExecStartPre=/bin/sleep 60
SuccessExitStatus=SIGUSR1
