# systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
OnFailure=notify-fail@%N.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va
