# systemd/user/notify-fail@.service
[Unit]
Description=Notify service failure

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send -u critical "Systemd" "Service %i failed"
