# systemd/user/sync.service
[Unit]
Description=Sync to cloud

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-cloud
