# systemd/user/updates-check.service
[Unit]
Description=Check for updates
OnFailure=notify-fail@%N.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/updates-check
