# systemd/user/updates-check.service
[Unit]
Description=Check for updates

[Service]
Type=oneshot
ExecStart=%h/.local/bin/updates-check
